# otun.timer
[Unit]
Description=Run otun.sh at a specific time

[Timer]
OnCalendar=*-*-* 04:33:00
AccuracySec=1s

[Install]
WantedBy=timers.target
