# otun.timer
[Unit]
Description=Run otun.sh at a specific time

[Timer]
OnCalendar=*-*-* HOUR_MINUTE:00
Persistent=true

[Install]
WantedBy=timers.target
