# otun.service
[Unit]
Description=otun (On-Telegram Updates Notifier)

[Service]
Type=simple
ExecStart=/opt/otun/otun.sh
WorkingDirectory=/opt/otun
