# otun.service
[Unit]
Description=otun (On-Telegram Updates Notifier)

[Service]
Type=simple
ExecStart=%h/otun/otun.sh
WorkingDirectory=%h/otun
